# deployment/root_setup/templates/backup-db.timer.tmpl
[Unit]
Description=Daily SQLite backup timer for Kistentracker

[Timer]
# Run hourly.
OnCalendar=hourly
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target