# systemd/caddy.socket.in
[Socket]
$LISTEN_STREAMS
[Install]
WantedBy=sockets.target